# agents/aegis/systemd/aegis.service
[Unit]
Description=Aegis Agent - eBPF Policy Enforcement
Documentation=https://github.com/sgerhart/aegis_agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aegis-agent
ExecStart=/opt/aegis-agent/aegis-agent
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment variables
Environment=AGENT_HOST_ID=
Environment=AGENT_ID=
Environment=AEGIS_DRY_RUN=false
Environment=AEGIS_MTLS=true
Environment=AEGIS_TRUSTSTORE=/etc/aegis/truststore.json
Environment=AEGIS_VERBOSE=false
Environment=AGENT_REGISTRY_URL=http://192.168.1.166:8090
Environment=NATS_URL=nats://192.168.1.166:4222

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aegis-agent
ReadWritePaths=/sys/fs/bpf
ReadWritePaths=/etc/aegis
ReadWritePaths=/var/lib/aegis

# Capabilities required for eBPF
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target
